# share/sunos/manifest/noc-launcher.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="noc-launcher">
    <service name="application/noc-launcher" type="service" version="1">
        <create_default_instance enabled='false' />
        <single_instance />
        <dependency name='loopback'
            grouping='require_all'
            restart_on='error'
            type='service'>
                <service_fmri value='svc:/network/loopback:default'/>
        </dependency>
        <dependency name='physical'
            grouping='optional_all'
            restart_on='error'
            type='service'>
                <service_fmri value='svc:/network/physical:default'/>
        </dependency>
        <dependency name='multiuser-server'
            grouping='require_all'
            restart_on='error'
            type='service'>
                <service_fmri value='svc:/milestone/multi-user-server:default'/>
        </dependency>
        <method_context working_directory="/opt/noc">
            <method_credential user="root" group="root"/>
            <method_environment>
                <envvar name="PATH" value="/opt/csw/bin:/usr/bin:/bin"/>
            </method_environment>
        </method_context>
        <exec_method type="method" name="start"
            exec="%{noc-launcher/python} %{noc-launcher/prefix}/scripts/noc-launcher.py -c %{noc-launcher/config} %m"
            timeout_seconds='60' />
        <exec_method type="method" name="stop"
            exec="%{noc-launcher/python} %{noc-launcher/prefix}/scripts/noc-launcher.py -c %{noc-launcher/config} %m"
            timeout_seconds='60' />
        <exec_method type="method" name="refresh"
            exec="%{noc-launcher/python} %{noc-launcher/prefix}/scripts/noc-launcher.py -c %{noc-launcher/config} %m"
            timeout_seconds='60' />
        <property_group name="noc-launcher" type="framework">
            <stability value="Evolving"/>
            <propval name="python"      type="astring" value="/opt/csw/bin/python" />
            <propval name="prefix"      type="astring" value="/opt/noc" />
            <propval name="config"      type="astring" value="/opt/noc/etc/noc-launcher.conf" />
        </property_group>
        <template>
            <common_name>
                <loctext xml:lang='C'>
                    NOC Launcher
                </loctext>
            </common_name>
        </template>
    </service>
</service_bundle>