# ansible/roles/redis/templates/etc/systemd/system/redis.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/redis-server {{ redis_conf_path }}/redis.conf --daemonize no
