# ansible/roles/node/templates/etc/config/services/activator.conf.j2
[program:activator]
autostart = {{'svc-activator' in group_names}}
