# ansible/roles/node/templates/etc/config/services/graphitecollector.conf.j2
[program:graphitecollector]
autostart = {{'svc-graphitecollector' in group_names}}
