# ansible/roles/node/templates/etc/config/services/nm.conf.j2
[program:nm]
autostart = {{'svc-nm' in group_names}}
