# ansible/roles/node/templates/etc/config/services/ping.conf.j2
[program:ping]
autostart = {{'svc-ping' in group_names}}
