# ansible/roles/node/templates/etc/config/services/pmwriter.conf.j2
[program:pmwriter]
autostart = {{'svc-pmwriter' in group_names}}
