# ansible/roles/node/templates/etc/config/services/probe.conf.j2
[program:probe]
autostart = {{'svc-probe' in group_names}}
