# ansible/roles/node/templates/etc/config/services/probeconf.conf.j2
[program:probeconf]
autostart = {{'svc-probeconf' in group_names}}
