# ansible/roles/node/templates/etc/config/services/syslogcollector.conf.j2
[program:syslogcollector]
autostart = {{'svc-syslogcollector' in group_names}}
