# ansible/roles/node/templates/etc/config/services/trapcollector.conf.j2
[program:trapcollector]
autostart = {{'svc-trapcollector' in group_names}}
