# ansible/roles/node/templates/etc/config/services/web.conf.j2
[program:web]
autostart = {{'svc-web' in group_names}}
command = ./scripts/noc-web.py --listen={{ inventory_hostname }}:8000