# ansible/roles/nsqadmin/templates/etc/systemd/nsqadmin.service
[Unit]
Description=Admin interface Realtime distributed messaging (nsddmin)
After=network.target

[Service]
ExecStart={{noc_root}}/bin/nsqadmin -config {{noc_root}}/var/etc/nsq/nsqadmin.conf
Restart=always
User=nsq
Group=nsq

[Install]
WantedBy=multi-user.target